# etc/cron.d/10-netadmin
30 3 * * * root cd /opt/netadmin/netadmin && ../venv/bin/python manage.py updateMaxMind
