# etc/cron.d/netadmin-updateMaxMind
30 3 * * * netadmin cd /opt/netadmin/netadmin && ../venv/bin/python manage.py updateMaxMind
